# profile.d/paths
# Shared PATH management for both zsh and bash
# This file is sourced by both zprofile and bash_profile

add-path /opt/homebrew/bin
add-path "$HOME/Library/Application Support/Coursier/bin"
add-path $HOME/.local/bin
add-path $HOME/.cargo/bin
add-path $HOME/.lmstudio/bin
add-path $HOME/.codeium/windsurf/bin
add-path $HOME/.claude/local